# software/coral/camera.service
[Unit]
Description=Camera service

[Service]
ExecStart=sh /home/mendel/kidz/software/coral/run_noserver.sh

[Install]
WantedBy=multi-user.target
